# packaging/launchd/com.takezoh.credproxyd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.takezoh.credproxyd</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/credproxyd</string>
        <string>--config</string>
        <string>/Users/Shared/credproxyd/config.toml</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/credproxyd.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/credproxyd.log</string>
</dict>
</plist>
